# systemd/photos-backup.service
[Unit]
Description=Photos Backup

[Service]
Type=simple
Restart=always
ExecStart=/home/markr/work/personal/credfeto-photos/systemd/backup.photos
