# systemd/photos-backup.timer
[Unit]
Description=Photos Backup

[Timer]
OnCalendar=*-*-* *:*/5:00
Persistent=true
Unit=photos-backup.service

[Install]
WantedBy=multi-user.target
